# 04 apache-2.4.3
yum -y install openssl-devel

cd /tmp
本地上传httpd-2.4.3.tar.gz apr-1.4.6.tar.gz apr-util-1.5.1.tar.gz pcre-8.31.tar.gz
tar zxvf pcre-8.31.tar.gz
tar zxvf apr-1.4.6.tar.gz
tar zxvf apr-util-1.5.1.tar.gz 
tar zxvf httpd-2.4.3.tar.gz

cd pcre-8.31 
./configure --prefix=/usr/local/pcre
make && make install

cd /tmp
cp -r /tmp/apr-1.4.6 /tmp/httpd-2.4.3/srclib/apr
cp -r /tmp/apr-util-1.5.1 /tmp/httpd-2.4.3/srclib/apr-util
cd httpd-2.4.3
./configure --prefix=/usr/local/apache \
	--enable-dav \
	--enable-so \
	--enable-maintainer-mode \
	--enable-rewrite \
	--enable-deflate=shared \
	--enable-ssl=shared \
	--enable-expires=shared \
	--enable-headers=shared \
	--enable-static-support \
	--with-included-apr \
	--with-mpm=prefork \
	--enable-cache \
	--enable-file-cache \
	--with-pcre=/usr/local/pcre
make && make install
cp /usr/local/apache/bin/apachectl /etc/init.d/httpd
vi /etc/init.d/httpd
在#!/bin/sh 下面加上这两行
# chkconfig: 345 90 90
# description: Apache
:wq 
chkconfig --add httpd
chkconfig httpd on

vi /usr/local/apache/conf/httpd.conf
Listen 8000
ServerName localhost:8000
User www
Group www 
:wq

service httpd start
vi /etc/sysconfig/iptables
-A INPUT -m state --state NEW -m tcp -p tcp --dport 8000 -j ACCEPT
:wq 
service iptables restart 

# groupadd www
# useradd www -g www -s /sbin/nologin
# service httpd restart
# /usr/local/apache/bin/httpd -v
# /usr/local/apache/bin/httpd -M
至此apache安装结束。

#apache验证
cd /usr/local/apache/htdocs
vi index.php
<?php
phpinfo();
?>
:wq!
# chown -R www.www /usr/local/apache/htdocs
#配置apache以module方式使用php 

# vi /usr/local/apache/conf/httpd.conf 
<IfModule mime_module> 
    #此处添加php handler 
    AddType application/x-httpd-php .php .phtml
    AddType application/x-httpd-php-source .phps
    PHPIniDir "/etc" 
</IfModule>
DirectoryIndex index.html index.php

#禁止文件列表功能 
#Options Indexes FollowSymLinks 
Options FollowSymLinks

#允许配置重写 
#AllowOverride None 
AllowOverride All
:wq 

# service httpd restart
# 访问 http://192.168.203.80/index.php